# legacy/docker/novnc/desktop-realfix.conf
; === Piper MoveIt Docker — real(host-network) desktop fix program ===
; Loaded by the base's supervisord via /etc/supervisor/supervisord.conf `[include] files =
; /etc/supervisor/conf.d/*.conf`. The base regenerates conf.d/supervisord.conf (vnc+novnc) on
; every start but NEVER touches this file, so our fix survives restarts/reboots.
;
; Runs once at startup. For MODE=real it relocates the desktop to display :2 and noVNC to 6080
; (see desktop-realfix.sh for the why). For every other profile it is an immediate no-op.
; priority=1500 (> vnc/novnc default 999) so it runs LAST and is the sole owner of the desktop's
; final state — no race with supervisord's autostart of vnc/novnc. It stops them (clearing the
; FATAL they hit on :1/port 80), cleans stale state, and starts them on :2/6080 with a retry
; loop. autorestart=false → runs exactly once.

[program:desktop-realfix]
command=/desktop-realfix.sh
user=root
autostart=true
autorestart=false
startsecs=0
startretries=1
priority=1500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; MODE is injected as a container env var (compose `real` service sets MODE=real).
environment=MODE="%(ENV_MODE)s"
